# snapcast-server/rootfs/etc/services.d/snapcast-server/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start snapcast service
# ==============================================================================

# Get config values
spotify_name=$(bashio::config 'spotify_name')
spotify_username=$(bashio::config 'spotify_username')
spotify_password=$(bashio::config 'spotify_password')
log_level=$(bashio::config 'log_level')

# Start snapserver
bashio::log.info "Starting snapcast server..."
bashio::log.info "Setting log level to '${log_level}'"

CONFIG=$(cat << EOF
{
  "streams": {
    "streams": []
  },
  "buffer": 1000,
  "codec": "flac",
  "send_to_muted": false,
  "sampleformat": "48000:16:2",
  "spotify_name": "${spotify_name}",
  "log_level": "${log_level}",
  "spotify_username": "${spotify_username}",
  "spotify_password": "${spotify_password}"
}
EOF
)

bashio::log.debug "${CONFIG}"

# Start snapserver with librespot stream
snapserver \
  -s "librespot:///usr/bin/librespot?name=Spotify&devicename=${spotify_name}&username=${spotify_username}&password=${spotify_password}&disable_audio_cache=true"
